# bin/phpbot404.conf
env USER="www-data"
env SCRIPT_DIR="/var/www/private/IRC-Bot-404"
env LOG_FILE="/var/log/phpbot404.log"

description "PHP Bot 404 startup script"
author      "matejvelikonja"

start on startup
stop on shutdown

script
    exec /usr/bin/sudo -u $USER /usr/bin/php $SCRIPT_DIR/phpbot404.php >> $LOG_FILE
end script